# deploy/systemd/opspilot.service
[Unit]
Description=OpsPilot API server
Documentation=https://github.com/vicenteliu/OpsPilot
After=network.target
Wants=network-online.target

[Service]
Type=exec
User=opspilot
Group=opspilot
WorkingDirectory=/opt/opspilot

# Set OPSPILOT_HOME and load secrets from /etc/opspilot/env (optional).
Environment=OPSPILOT_HOME=/var/lib/opspilot
EnvironmentFile=-/etc/opspilot/env

ExecStart=/opt/opspilot/.venv/bin/opspilot serve \
    --host 127.0.0.1 \
    --port 8000 \
    --workers 2 \
    --json-logs

# Restart policy
Restart=always
RestartSec=5s
TimeoutStopSec=30s

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/opspilot
PrivateTmp=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
